# deploy/skipperbot-deploy-watcher.service
# Optional host-side deploy watcher for the UI/API "deploy + restart" action.
# Runs scripts/deploy_watcher.sh on the host: when the agent drops a
# .deploy_pending sentinel (after a graceful drain), this pulls the latest code
# and recycles the docker compose stack. Without it, the deploy button is just a
# plain container restart (no code update).
#
# Install on the host (adjust User/WorkingDirectory to your checkout):
#   sudo cp deploy/skipperbot-deploy-watcher.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now skipperbot-deploy-watcher
#
# Requires: the same user can run `git pull` in the repo and `docker compose`.

[Unit]
Description=Skipperbot deploy watcher (git pull + docker recycle on UI/API deploy)
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=rodney
WorkingDirectory=/home/rodney/repos/skipperbot-platform
ExecStart=/home/rodney/repos/skipperbot-platform/scripts/deploy_watcher.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
